# board/diasom/rk3588_b2b_evb/fs-overlay/usr/lib/systemd/system/rtc-init.service
[Unit]
Description=Initialize RTC without battery
DefaultDependencies=no
After=dev-rtc0.device
Before=time-sync.target

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'hwclock -w'

[Install]
WantedBy=multi-user.target